# supervisor/conf.d/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=077

[supervisord]
nodaemon=true
logfile=/apps/logs/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# 这个集群族节点的配置
[program:openresty]
command=/bin/bash -c "source /sh/init_openresty.sh"
stderr_logfile=/apps/logs/supervisor_%(program_name)s_stderr.log
stdout_logfile=/apps/logs/supervisor_%(program_name)s_stdout.log

[program:user-rpc-server]
command=/bin/bash -c "cd /apps/go_user_rpc && go run main.go";
stderr_logfile=/apps/logs/supervisor_%(program_name)s_stderr.log
stdout_logfile=/apps/logs/supervisor_%(program_name)s_stdout.log

[program:user-server]
command=/sh/wait-for-it.sh 127.0.0.1:8080 -- /bin/bash -c "cd /apps/go_user && go run main.go";
stderr_logfile=/apps/logs/supervisor_%(program_name)s_stderr.log
stdout_logfile=/apps/logs/supervisor_%(program_name)s_stdout.log

[supervisorctl]              ;这个主要是针对supervisorctl的一些配置
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

